# scripts/com.wrc.quiz-weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template — do not load directly. Replace placeholders before installing:
    __REPO_DIR__  absolute path to the permanent clone of this repo
    __HOME__      absolute path to the owning user's home directory

  See docs/quiz-cron.md for the install/uninstall commands (sed substitution
  + launchctl bootstrap/bootout). This launches quiz-weekly-cron.sh, which
  only collects real data, generates a draft, and stops at the human
  approval queue — it never publishes anything on its own.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.wrc.quiz-weekly</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>__REPO_DIR__/scripts/quiz-weekly-cron.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/quiz-weekly.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/quiz-weekly.log</string>
</dict>
</plist>
